# codetest.service
[Unit]
Description=RetroFit Image Capture Service v2.1
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/Desktop/Evaratech/Evaraflow

# Python environment
Environment="PYTHONDONTWRITEBYTECODE=1"
Environment="PYTHONUNBUFFERED=1"
Environment="VIRTUAL_ENV=/home/pi/Desktop/Evaratech/Evaraflow/.venv"
Environment="PATH=/home/pi/Desktop/Evaratech/Evaraflow/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Pre-start: wait for network + clear stale .pyc files
ExecStartPre=/bin/sleep 10
ExecStartPre=/usr/bin/find /home/pi/Desktop/Evaratech/Evaraflow -type f -name "*.pyc" -delete
ExecStartPre=/usr/bin/find /home/pi/Desktop/Evaratech/Evaraflow -type d -name "__pycache__" -exec rm -rf {} +

# Start using venv Python
ExecStart=/home/pi/Desktop/Evaratech/Evaraflow/.venv/bin/python3 /home/pi/Desktop/Evaratech/Evaraflow/main_service.py

# Restart policy
Restart=always
RestartSec=30

# Resource limits (safety for Pi Zero W)
MemoryMax=200M
CPUQuota=90%

[Install]
WantedBy=multi-user.target
